# assets/supervisor/explorer.conf
[program:explorer]
priority=1020
command=wine start /max /d "/root/.wine/drive_c/windows/system32/" explorer.exe