# assets/supervisor/fluxbox.conf
[program:fluxbox]
user=vncer
priority=1020
command=fluxbox -display :0
autorestart=false
