# assets/supervisor/vncserver.conf
[program:vncserver]
priority=1020
command=vncserver :0 -passwordfile /home/winer/.vnc/passwd -rfbport 5900
autorestart=false