# assets/supervisor/vncserver.conf
[program:vncserver]
user=vncer
priority=1020
command=vncserver :0 -passwordfile /home/vncer/.vnc/passwd -rfbport 5900
autorestart=false
