# assets/supervisor/vncserver.conf
[program:vncserver]
user=vncer
priority=1020
command=/usr/bin/startvncserver.sh
autorestart=false
