# assets/supervisor/xfce4.conf
[program:xfce4]
priority=1020
command=startxfce4 
autorestart=false
