# assets/supervisor/xvfb.conf
[program:xvfb]
priority=1000
command=Xvfb :0 -ac -screen 0 1024x768x16
autorestart=false